# config/unison/unison-usb.service
[Unit]
Description=Unison backup-usb

[Service]
ExecStart=unison backup-usb
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
